# Engram-AiCrawler/01_devroot/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/data/logs/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/app/data/logs
loglevel=info


[program:crawl4ai]
command=/usr/local/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 11235
directory=/app

autostart=true
autorestart=true
startretries=3
startsecs=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/app/data/logs/crawl4ai-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/app/data/logs/crawl4ai-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONUNBUFFERED="1",APP_NAME="%(ENV_APP_NAME)s",APP_VERSION="%(ENV_APP_VERSION)s",DEBUG="%(ENV_DEBUG)s",LOG_LEVEL="%(ENV_LOG_LEVEL)s",API_HOST="%(ENV_API_HOST)s",API_PORT="%(ENV_API_PORT)s",CORS_ORIGINS="%(ENV_CORS_ORIGINS)s",CRAWL4AI_VERSION="%(ENV_CRAWL4AI_VERSION)s",BROWSER_TYPE="%(ENV_BROWSER_TYPE)s",HEADLESS="%(ENV_HEADLESS)s",VIEWPORT_WIDTH="%(ENV_VIEWPORT_WIDTH)s",VIEWPORT_HEIGHT="%(ENV_VIEWPORT_HEIGHT)s",CACHE_ENABLED="%(ENV_CACHE_ENABLED)s",CACHE_MODE="%(ENV_CACHE_MODE)s",WORD_COUNT_THRESHOLD="%(ENV_WORD_COUNT_THRESHOLD)s",PAGE_TIMEOUT="%(ENV_PAGE_TIMEOUT)s",LM_STUDIO_URL="%(ENV_LM_STUDIO_URL)s",LM_STUDIO_MODEL="%(ENV_LM_STUDIO_MODEL)s",LM_STUDIO_TIMEOUT="%(ENV_LM_STUDIO_TIMEOUT)s",LM_STUDIO_MAX_RETRIES="%(ENV_LM_STUDIO_MAX_RETRIES)s",LM_STUDIO_RETRY_DELAY="%(ENV_LM_STUDIO_RETRY_DELAY)s",LM_STUDIO_TEMPERATURE="%(ENV_LM_STUDIO_TEMPERATURE)s",REDIS_URL="%(ENV_REDIS_URL)s",REDIS_MAX_CONNECTIONS="%(ENV_REDIS_MAX_CONNECTIONS)s",REDIS_SOCKET_TIMEOUT="%(ENV_REDIS_SOCKET_TIMEOUT)s",REDIS_SOCKET_CONNECT_TIMEOUT="%(ENV_REDIS_SOCKET_CONNECT_TIMEOUT)s",CHROMADB_PATH="%(ENV_CHROMADB_PATH)s",CHROMADB_COLLECTION_PREFIX="%(ENV_CHROMADB_COLLECTION_PREFIX)s",CHROMADB_SIMILARITY_THRESHOLD="%(ENV_CHROMADB_SIMILARITY_THRESHOLD)s",CHROMADB_EMBEDDING_BATCH_SIZE="%(ENV_CHROMADB_EMBEDDING_BATCH_SIZE)s",DATA_HOT_PATH="%(ENV_DATA_HOT_PATH)s",DATA_WARM_PATH="%(ENV_DATA_WARM_PATH)s",DATA_COLD_PATH="%(ENV_DATA_COLD_PATH)s",DATA_ARCHIVE_PATH="%(ENV_DATA_ARCHIVE_PATH)s",WATCHDOG_CHECK_INTERVAL_SECONDS="%(ENV_WATCHDOG_CHECK_INTERVAL_SECONDS)s",WATCHDOG_ORPHAN_AGE_MINUTES="%(ENV_WATCHDOG_ORPHAN_AGE_MINUTES)s",WATCHDOG_MEMORY_THRESHOLD_PERCENT="%(ENV_WATCHDOG_MEMORY_THRESHOLD_PERCENT)s",WATCHDOG_DISK_THRESHOLD_PERCENT="%(ENV_WATCHDOG_DISK_THRESHOLD_PERCENT)s"

[program:lm_bridge]
command=/usr/local/bin/python -m app.services.lm_studio_bridge
directory=/app

autostart=false
autorestart=false
startretries=3
startsecs=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/app/data/logs/lm_bridge-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/app/data/logs/lm_bridge-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONUNBUFFERED="1",APP_NAME="%(ENV_APP_NAME)s",APP_VERSION="%(ENV_APP_VERSION)s",DEBUG="%(ENV_DEBUG)s",LOG_LEVEL="%(ENV_LOG_LEVEL)s",LM_STUDIO_URL="%(ENV_LM_STUDIO_URL)s",LM_STUDIO_MODEL="%(ENV_LM_STUDIO_MODEL)s",LM_STUDIO_TIMEOUT="%(ENV_LM_STUDIO_TIMEOUT)s",LM_STUDIO_MAX_RETRIES="%(ENV_LM_STUDIO_MAX_RETRIES)s",LM_STUDIO_RETRY_DELAY="%(ENV_LM_STUDIO_RETRY_DELAY)s",LM_STUDIO_TEMPERATURE="%(ENV_LM_STUDIO_TEMPERATURE)s",CHROMADB_PATH="%(ENV_CHROMADB_PATH)s",CHROMADB_COLLECTION_PREFIX="%(ENV_CHROMADB_COLLECTION_PREFIX)s",CHROMADB_SIMILARITY_THRESHOLD="%(ENV_CHROMADB_SIMILARITY_THRESHOLD)s",CHROMADB_EMBEDDING_BATCH_SIZE="%(ENV_CHROMADB_EMBEDDING_BATCH_SIZE)s",WATCHDOG_CHECK_INTERVAL_SECONDS="%(ENV_WATCHDOG_CHECK_INTERVAL_SECONDS)s",WATCHDOG_ORPHAN_AGE_MINUTES="%(ENV_WATCHDOG_ORPHAN_AGE_MINUTES)s",WATCHDOG_MEMORY_THRESHOLD_PERCENT="%(ENV_WATCHDOG_MEMORY_THRESHOLD_PERCENT)s",WATCHDOG_DISK_THRESHOLD_PERCENT="%(ENV_WATCHDOG_DISK_THRESHOLD_PERCENT)s"

[program:watchdog]
command=/usr/local/bin/python -m app.services.watchdog
directory=/app

autostart=true
autorestart=true
startretries=3
startsecs=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/app/data/logs/watchdog-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/app/data/logs/watchdog-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONUNBUFFERED="1",APP_NAME="%(ENV_APP_NAME)s",APP_VERSION="%(ENV_APP_VERSION)s",DEBUG="%(ENV_DEBUG)s",LOG_LEVEL="%(ENV_LOG_LEVEL)s",WATCHDOG_CHECK_INTERVAL_SECONDS="%(ENV_WATCHDOG_CHECK_INTERVAL_SECONDS)s",WATCHDOG_ORPHAN_AGE_MINUTES="%(ENV_WATCHDOG_ORPHAN_AGE_MINUTES)s",WATCHDOG_MEMORY_THRESHOLD_PERCENT="%(ENV_WATCHDOG_MEMORY_THRESHOLD_PERCENT)s",WATCHDOG_DISK_THRESHOLD_PERCENT="%(ENV_WATCHDOG_DISK_THRESHOLD_PERCENT)s",WATCHDOG_MEMORY_CHECK_ENABLED="%(ENV_WATCHDOG_MEMORY_CHECK_ENABLED)s",WATCHDOG_DISK_CHECK_ENABLED="%(ENV_WATCHDOG_DISK_CHECK_ENABLED)s",CHROMADB_PATH="%(ENV_CHROMADB_PATH)s"

[program:cleanup]
command=/usr/local/bin/python -m app.services.data_lifecycle
directory=/app

autostart=true
autorestart=true
startretries=3
startsecs=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/app/data/logs/cleanup-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/app/data/logs/cleanup-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONUNBUFFERED="1",APP_NAME="%(ENV_APP_NAME)s",APP_VERSION="%(ENV_APP_VERSION)s",DEBUG="%(ENV_DEBUG)s",LOG_LEVEL="%(ENV_LOG_LEVEL)s",DATA_HOT_PATH="%(ENV_DATA_HOT_PATH)s",DATA_WARM_PATH="%(ENV_DATA_WARM_PATH)s",DATA_COLD_PATH="%(ENV_DATA_COLD_PATH)s",DATA_ARCHIVE_PATH="%(ENV_DATA_ARCHIVE_PATH)s",DATA_HOT_MAX_AGE_HOURS="%(ENV_DATA_HOT_MAX_AGE_HOURS)s",DATA_WARM_MAX_AGE_DAYS="%(ENV_DATA_WARM_MAX_AGE_DAYS)s",DATA_COLD_MAX_AGE_DAYS="%(ENV_DATA_COLD_MAX_AGE_DAYS)s",DATA_OFFLOAD_THRESHOLD_DAYS="%(ENV_DATA_OFFLOAD_THRESHOLD_DAYS)s",DATA_ARCHIVE_THRESHOLD_GB="%(ENV_DATA_ARCHIVE_THRESHOLD_GB)s",DATA_CLEANUP_INTERVAL_MINUTES="%(ENV_DATA_CLEANUP_INTERVAL_MINUTES)s",CHROMADB_PATH="%(ENV_CHROMADB_PATH)s",WATCHDOG_DISK_THRESHOLD_PERCENT="%(ENV_WATCHDOG_DISK_THRESHOLD_PERCENT)s"

[program:arq_worker]
command=/usr/local/bin/python -m arq app.workers.arq_worker.WorkerSettings
directory=/app

autostart=true
autorestart=true
startretries=3
startsecs=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/app/data/logs/arq_worker-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/app/data/logs/arq_worker-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONUNBUFFERED="1",REDIS_URL="%(ENV_REDIS_URL)s",LOG_LEVEL="%(ENV_LOG_LEVEL)s",HEADLESS="%(ENV_HEADLESS)s",BROWSER_TYPE="%(ENV_BROWSER_TYPE)s"

[group:services]
programs=crawl4ai,watchdog,cleanup,arq_worker
priority=999
